# cron/reboot_on_chromium_fail.cron
*/2 * * * * /home/pi/scripts/reboot_on_chromium_fail
