# cron/reboot_on_browser_fail.cron
*/3 * * * * /usr/local/bin/reboot_on_browser_fail
