# cron/refresh_browser.cron
*/5 * * * * /usr/local/bin/refresh_browser
