# cron/reset_network_on_fail.cron
*/2 * * * * /usr/local/bin/reset_network_on_fail
